# scripts/ai.openclaw.leadflow.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.leadflow.heartbeat</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>heartbeat-wrapper.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/clawdbot/.openclaw/genome/core</string>
    
    <key>StartInterval</key>
    <integer>1800</integer>
    <!-- Run every 30 minutes (1800 seconds) — hot path handled by realtime-dispatcher -->
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/openclaw/leadflow-heartbeat.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/leadflow-heartbeat-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
